# supervisor.conf
[program:graphite]
command=/sbin/my_init

[program:save_env_for_cron]
command=/bin/bash -c 'export | grep DELETE_OLD_FILES > /opt/env.saved'
autorestart=false
startsecs=0
